# config/systemd/user/wmf-fingerprints.timer
[Unit]
Description=Run wmf-fingerprints daily

[Timer]
OnBootSec=30min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
